# image/etc/supervisor/conf.d/supervisord.conf
[program:sshd]
priority=20
directory=/usr/bin/sshd
command=/usr/bin/sshd -D
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/sshd.log
redirect_stderr=true

